# roles/dummy/files/dummy.service
[Unit]
Description=Dummy

[Service]
ExecStart=/bin/sleep infinity

[Install]
WantedBy=default.target